# scripts/launchd/com.cs2pricer.collector.plist
<?xml version="1.0" encoding="UTF-8"?>
<!--
  macOS LaunchAgent for the sold-price collector. Runs scripts/collect.py every 12h.
  Paths are absolute to this machine's checkout; adjust if you move the repo.

  Install:
    cp scripts/launchd/com.cs2pricer.collector.plist ~/Library/LaunchAgents/
    launchctl load ~/Library/LaunchAgents/com.cs2pricer.collector.plist
  Status / logs:
    launchctl list | grep cs2pricer
    tail -f data/collector/collect.log
  Stop:
    launchctl unload ~/Library/LaunchAgents/com.cs2pricer.collector.plist
-->
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.cs2pricer.collector</string>

    <key>ProgramArguments</key>
    <array>
        <string>/Users/dna/coding-projects/cs2-skin-pricer/.venv/bin/python</string>
        <string>/Users/dna/coding-projects/cs2-skin-pricer/scripts/collect.py</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/dna/coding-projects/cs2-skin-pricer</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PYTHONPATH</key>
        <string>/Users/dna/coding-projects/cs2-skin-pricer</string>
    </dict>

    <!-- Every 12 hours. Each run is a ~100-request burst, well under the 200/window cap. -->
    <key>StartInterval</key>
    <integer>43200</integer>

    <key>RunAtLoad</key>
    <true/>

    <key>StandardOutPath</key>
    <string>/Users/dna/coding-projects/cs2-skin-pricer/data/collector/collect.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/dna/coding-projects/cs2-skin-pricer/data/collector/collect.log</string>
</dict>
</plist>
