# scripts/launchd/com.cs2pricer.backup.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!-- Daily pull-backup of the Pi collector DB. StartCalendarInterval (unlike
     StartInterval) fires a missed slot once on wake, so the Mac sleeping
     through 10:00 just delays the backup instead of skipping the day.
     Install: cp to ~/Library/LaunchAgents/ && launchctl load it. -->
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.cs2pricer.backup</string>
    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>/Users/dna/coding-projects/cs2-skin-pricer/scripts/backup_pull.sh</string>
    </array>
    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>10</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>
    <key>StandardOutPath</key>
    <string>/Users/dna/coding-projects/cs2-skin-pricer/data/collector/backup_pull.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/dna/coding-projects/cs2-skin-pricer/data/collector/backup_pull.log</string>
</dict>
</plist>
